# .config/autostart/solaar.desktop
[Desktop Entry]
Name=Solaar
Type=Application
Exec=/usr/bin/solaar -w hide -b solaar 
X-GNOME-Autostart-enabled=true
