# Documents/backup/etc/crontab
# m h dom mon dow command
@hourly /home/$USER/.local/statusbar/sb-clock
@daily /home/$USER/.local/statusbar/moonphase
@hourly /home/$USER/.local/statusbar/openweather-emoji
